# ops/rocketchat/templates/com.velocityworks.ngrok-rocketchat.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- IMP-11: rendered by install-launchd.sh -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.velocityworks.ngrok-rocketchat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@ROOT@/run_ngrok.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>@HOME@/logs/ngrok-rocketchat/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>@HOME@/logs/ngrok-rocketchat/launchd.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>@HOME@</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
